# importantes/crontab
0 0 * * 1 >/var/log/svxlink.log
0 0 * * 1 >/var/log/svxreflector.log
